# openstack/nova-network.service
[Unit]
Description=Nova Network Service

[Service]
ExecStart=/opt/openstack/bin/nova-network
Restart=always

[Install]
WantedBy=multi-user.target
