# openstack/nova-volume.service
[Unit]
Description=Nova API Service

[Service]
ExecStart=/opt/openstack/bin/nova-volume
Restart=always

[Install]
WantedBy=multi-user.target
