# openstack-git/glance-api.service
[Unit]
Description=Glance API Service

[Service]
ExecStart=/opt/openstack-essex/bin/glance-api
Restart=always

[Install]
WantedBy=multi-user.target
