# openstack-git/keystone.service
[Unit]
Description=Keystone Identity Service

[Service]
ExecStart=/opt/openstack-essex/bin/keystone-all
Restart=always

[Install]
WantedBy=multi-user.target
