# openstack-git/nova-api.service
[Unit]
Description=Nova API Service

[Service]
ExecStart=/opt/openstack/bin/nova-api
Restart=always

[Install]
WantedBy=multi-user.target
